# docker/supervisord.conf.d/webui.conf
[program:webui]
directory=/app
priority=100
command = pdm run webui
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /app/logs/webui.log
stdout_logfile_maxbytes = 10MB

